# scripts/systemd/roon-worker.service
[Unit]
Description=Roon Helper Worker Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/RoonServerLHD
ExecStart=/usr/bin/env python -m worker.main
Restart=on-failure
EnvironmentFile=-/etc/roonhelper.env

[Install]
WantedBy=multi-user.target
